# backup.rules
SUBSYSTEMS=="block", ENV{ID_FS_UUID}=="UUID GOES HERE", RUN+="/etc/scripts/backup.sh"
